# userland/etc/mycelium/boot.target
[Unit]
Description=early boot synchronization point
Wants=multi-user.target

[Install]
WantedBy=multi-user.target
