# userland/etc/mycelium/multi-user.target
[Unit]
Description=normal interactive multi-user boot
Wants=console.service

[Install]
WantedBy=boot.target
