# userland/etc/mycelium/shutdown.target
[Unit]
Description=orderly shutdown synchronization point
Conflicts=multi-user.target
